# conf.d/fingui.conf
[program:fingui]
command=xterm -e /app/finish-gui.sh
autorestart=false
